# scripts/.service
[Unit]
Description = gostart-crm server
ConditionPathExists=/home/www/app/gostart-crm/server
After = network.target
StartLimitIntervalSec=500
StartLimitBurst=5

[Service]
User = www
Group = www-data
PermissionsStartOnly = true

WorkingDirectory = /home/www/app/gostart-crm
ExecStart = /home/www/app/gostart-crm/server
ExecStartPre = /bin/mkdir -pv /var/run/gostart-crm /var/log/gostart-crm
ExecStartPre = /bin/chown teh:teh /var/run/gostart-crm /var/log/gostart-crm

StandardOutput=append:/var/log/gostart-crm/stdout.log
StandardError=append:/var/log/gostart-crm/stderr.log

KillSignal=SIGINT
Restart=on-failure
RestartSec=5s

[Install]
WantedBy = multi-user.target
